# etc/systemd/system/getty@tty1.service.d/autologin-on-boot.conf
[Service]
ExecStart=
ExecStart=-/bin/autologin-on-boot %I
