# etc/systemd/system/getty@tty1.service.d/autologin-on-boot.conf
[Service]
ExecStart=
ExecStart=-/usr/share/autologin-on-boot/getty %I
